# etc/book.oscar-system.org.path
[Unit]
Description="Check whether OSCAR webhook requested update"

[Path]
PathChanged=/srv/www/www-oscar-book/data/book.oscar-system.org.trigger
Unit=book.oscar-system.org.service

[Install]
WantedBy=default.target
